# debian/irgsh/DEBIAN/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create irgsh user and group if they don't exist
        if ! getent group irgsh >/dev/null; then
            addgroup --system irgsh
        fi

        if ! getent passwd irgsh >/dev/null; then
            adduser --system --home /var/lib/irgsh --no-create-home \
                    --ingroup irgsh --disabled-password --shell /bin/bash \
                    --gecos "IRGSH System User" irgsh
        fi

        # Ensure home directory exists
        if [ ! -d /var/lib/irgsh ]; then
            mkdir -p /var/lib/irgsh
        fi

        # Set proper ownership and permissions
        chown -R irgsh:irgsh /var/lib/irgsh
        chown -R irgsh:irgsh /var/log/irgsh
        chmod 755 /var/lib/irgsh
        chmod 755 /var/log/irgsh

        # Make config directory readable by irgsh user
        chgrp irgsh /etc/irgsh
        chmod 750 /etc/irgsh

        # Config file should be readable by irgsh
        if [ -f /etc/irgsh/config.yaml ]; then
            chown root:irgsh /etc/irgsh/config.yaml
            chmod 640 /etc/irgsh/config.yaml
        fi

        # Add irgsh user to docker group if docker is installed
        if getent group docker >/dev/null; then
            usermod -aG docker irgsh || true
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# Automatically added by dh_installsystemd/13.14.1ubuntu5
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'irgsh-builder.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'irgsh-builder.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'irgsh-builder.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'irgsh-builder.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'irgsh-builder.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.14.1ubuntu5
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'irgsh-chief.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'irgsh-chief.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'irgsh-chief.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'irgsh-chief.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'irgsh-chief.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.14.1ubuntu5
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'irgsh-iso.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'irgsh-iso.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'irgsh-iso.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'irgsh-iso.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'irgsh-iso.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.14.1ubuntu5
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'irgsh-repo.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'irgsh-repo.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'irgsh-repo.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'irgsh-repo.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'irgsh-repo.service' >/dev/null || true
fi
# End automatically added section


exit 0
